# crontab
@reboot /path/to/check_updates.fish
0 * * * * /path/to/check_updates.fish
